# launchd/com.flow.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.flow.server</string>

    <!-- Call venv python directly in FOREGROUND — do NOT use start_flow_server.sh.
         That script uses nohup + background (&), which causes launchd to lose the
         process PID immediately and restart the wrapper every ThrottleInterval.
         launchd requires a process that stays in the foreground.
         HOMEPLACEHOLDER is substituted with $HOME by install_launchd.sh. -->
    <key>ProgramArguments</key>
    <array>
        <string>HOMEPLACEHOLDER/ofa_jack_agent/projects/flow/.venv/bin/python3</string>
        <string>server_local.py</string>
    </array>

    <!-- WorkingDirectory must be set so server_local.py can find its siblings
         (models, config files). cd to the repo root. -->
    <key>WorkingDirectory</key>
    <string>HOMEPLACEHOLDER/ofa_jack_agent/projects/flow</string>

    <!-- Venv python needs HOME for mlx cache, Piper model paths, etc. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>HOMEPLACEHOLDER</string>
        <key>PATH</key>
        <string>HOMEPLACEHOLDER/ofa_jack_agent/projects/flow/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on any non-zero exit or crash. Do NOT restart on clean exit (0)
         so a deliberate "stop server" is honoured without fighting launchd. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 30s minimum between restarts — covers mlx-whisper + Piper model warmup
         so a crash loop doesn't hammer the GPU during recovery. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Logs on internal SSD — always writable, survives any future volume changes.
         HOMEPLACEHOLDER → real $HOME via install_launchd.sh. -->
    <key>StandardOutPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_server.log</string>
    <key>StandardErrorPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_server.log</string>

</dict>
</plist>
