# launchd/com.flow.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.flow.tunnel</string>

    <!-- cloudflared binary on internal SSD (Homebrew) — no volume dependency.
         Reads tunnel config from ~/.cloudflared/config.yml (home dir, always available).
         Path must match: which cloudflared on the Mac mini. -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>flow</string>
    </array>

    <!-- No WorkingDirectory — cloudflared reads ~/.cloudflared/config.yml directly -->

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <!-- Home-dir log — no volume dependency.
         HOMEPLACEHOLDER → real $HOME via install_launchd.sh. -->
    <key>StandardOutPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_tunnel.log</string>

</dict>
</plist>
