# launchd/com.flow.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.flow.tunnel</string>

    <!-- cloudflared binary on internal SSD (Homebrew) — no volume dependency.
         Reads tunnel credentials from ~/.cloudflared/<tunnel-uuid>.json.
         PREREQUISITE: cloudflared tunnel login must have been run first.
         Path verified: which cloudflared → /opt/homebrew/bin/cloudflared. -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>flow</string>
    </array>

    <!-- cloudflared resolves credential paths relative to $HOME (~/.cloudflared/).
         Explicit HOME ensures the correct user home is visible in a launchd context. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>HOMEPLACEHOLDER</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash or non-zero exit. Do NOT restart on clean exit (0)
         so a deliberate kill is honoured without fighting launchd. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 15s between restarts — tunnel reconnects faster than server coldstart. -->
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <!-- Logs on internal SSD. HOMEPLACEHOLDER → real $HOME via install_launchd.sh. -->
    <key>StandardOutPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_tunnel.log</string>

</dict>
</plist>
