# launchd/com.flow.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.flow.server</string>

    <!-- Points to ~/.flow/launch_server.sh — internal SSD, always available.
         The wrapper waits for /Volumes/homelab to mount, then delegates to
         the canonical start_flow_server.sh on the external volume.
         HOMEPLACEHOLDER is substituted with $HOME by install_launchd.sh. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOMEPLACEHOLDER/.flow/launch_server.sh</string>
    </array>

    <!-- No WorkingDirectory — launch_server.sh cds into the repo after
         the volume is confirmed mounted. start_flow_server.sh also cds. -->

    <!-- Homebrew python3 must be on PATH for the server to start -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Home-dir log — must be on internal SSD so launchd can open the fd
         before /Volumes/homelab is mounted. HOMEPLACEHOLDER → real $HOME. -->
    <key>StandardOutPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_server.log</string>
    <key>StandardErrorPath</key>
    <string>HOMEPLACEHOLDER/.flow/launchd_server.log</string>

</dict>
</plist>
